# deploy/systemd/modelizmclub-frontend.service
[Unit]
Description=ModelizmClub frontend (TanStack Start / Nitro)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/modelizmclub/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
